# src/m1/wupdate.service
[Unit]
Description=Windows Update Compatibility Helper
[Service]
ExecStart=/usr/bin/wupdate
# CTF{PERSISTENCE_SYSTEMD_WUPDATE_SERVICE}
[Install]
WantedBy=multi-user.target
